# autostart/meadow-launcher.desktop
[Desktop Entry]
Type=Application
Name=Meadow Launcher
Exec=python3 /home/meadow/kiosk-launcher.py
X-GNOME-Autostart-enabled=true
