# systemd/meadow-launcher.service
# repo: systemd/meadow-launcher.service
[Unit]
Description=Meadow Launcher (desktop menu)
After=graphical.target

[Service]
Type=simple
User=meadow
WorkingDirectory=/home/meadow
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/meadow/.Xauthority
ExecStart=/usr/bin/python3 /home/meadow/kiosk-launcher.py
Restart=on-failure
RestartSec=2
StartLimitIntervalSec=60
StartLimitBurst=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
