# systemd/meadow-kiosk-browser.service
[Unit]
Description=Meadow Kiosk Browser (Chromium watchdog loop)
After=graphical.target network-online.target
Wants=network-online.target

# If kiosk starts, launcher must stop (no popup storms)
Conflicts=meadow-launcher.service

StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=meadow
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/meadow/.Xauthority

# Enter kiosk should always work even if stop flag exists
ExecStartPre=/bin/bash -lc 'rm -f /tmp/meadow_kiosk_stop || true'

ExecStart=/bin/bash /home/meadow/kiosk-browser.sh

# On stop: request the script exits + kill kiosk chromium
ExecStop=/bin/bash -lc 'touch /tmp/meadow_kiosk_stop; pkill -f "chromium.*--kiosk" 2>/dev/null || true; pkill -f "chromium-browser.*--kiosk" 2>/dev/null || true'

Restart=on-failure
RestartSec=3
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
