# systemd/meadow-launcher.service
[Unit]
Description=Meadow Launcher (desktop menu)
After=graphical.target

# If launcher starts, kiosk browser must stop
Conflicts=meadow-kiosk-browser.service

StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=meadow
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/meadow/.Xauthority

ExecStart=/usr/bin/python3 /home/meadow/kiosk-launcher.py
Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
