# systemd/meadow-kiosk.service
[Unit]
Description=Meadow Kiosk Service
After=network-online.target systemd-resolved.service
Wants=network-online.target systemd-resolved.service

StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=meadow
WorkingDirectory=/home/meadow/meadow-kiosk

# Ensure /run/meadow exists on every start (tmpfs) and is owned by meadow
RuntimeDirectory=meadow
RuntimeDirectoryMode=0755

Environment=PYTHONUNBUFFERED=1
Environment=GIT_TERMINAL_PROMPT=0

Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/meadow/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000

Environment=MEADOW_ROOT=/home/meadow/meadow-kiosk
Environment=MEADOW_STATE_DIR=/home/meadow/meadow-kiosk/state
Environment=MEADOW_KIOSK_URL_FILE=/home/meadow/meadow-kiosk/kiosk.url
Environment=MEADOW_OFFLINE_URL=file:///home/meadow/meadow-kiosk/offline.html

Environment=MEADOW_DAEMON_HEARTBEAT_FILE=/run/meadow/daemon_heartbeat
Environment=MEADOW_HEARTBEAT_MODE=legacy

Environment=MEADOW_KIOSK_STOP_FLAG=/run/meadow/kiosk_stop
Environment=MEADOW_UI_HEARTBEAT_FILE=/run/meadow/ui_heartbeat
Environment=MEADOW_WP_HEARTBEAT_FILE=/run/meadow/wp_heartbeat
Environment=MEADOW_KIOSK_PIDFILE=/run/meadow/kiosk_browser.pid
Environment=MEADOW_RESTART_LOG=/run/meadow/kiosk_restart_times

Environment=MEADOW_KILL_POPUPS=0
Environment=MEADOW_KILL_POPUP_INTERVAL=2

ExecStartPre=/bin/bash -lc 'cd /home/meadow/meadow-kiosk && git fetch --prune origin'
ExecStartPre=/bin/bash -lc 'cd /home/meadow/meadow-kiosk && git reset --hard origin/main'
ExecStart=/usr/bin/python3 /home/meadow/meadow-kiosk/pi_api.py

Restart=on-failure
RestartSec=5

# Keep some safe hardening, but DO NOT break Chromium namespace/sandbox
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
UMask=007

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
