# systemd/meadow-kiosk.service
[Unit]
Description=Meadow Kiosk Service
After=network-online.target systemd-resolved.service
Wants=network-online.target systemd-resolved.service

[Service]
Type=simple
User=meadow
WorkingDirectory=/home/meadow/meadow-kiosk
Environment=PYTHONUNBUFFERED=1
Environment=GIT_TERMINAL_PROMPT=0
ExecStartPre=/bin/bash -lc 'cd /home/meadow/meadow-kiosk && git fetch --prune origin'
ExecStartPre=/bin/bash -lc 'cd /home/meadow/meadow-kiosk && git reset --hard origin/main'
# Local HTTP API for the kiosk UI -> Pi (Sigma + vend)
ExecStart=/usr/bin/python3 /home/meadow/meadow-kiosk/pi_api.py

# Resilience
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=10

# Hardening (safe defaults for this app)
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictNamespaces=true
UMask=007

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
