# systemd/meadow-kiosk-freeze-watchdog.service
[Unit]
Description=Meadow Kiosk Freeze Watchdog (restart Chromium if screen frozen)
After=graphical.target meadow-kiosk.service
Wants=meadow-kiosk.service
PartOf=meadow-kiosk.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=meadow
Environment=XDG_RUNTIME_DIR=/run/user/1000
Restart=always
RestartSec=5

# You can override these via systemd drop-in if needed
Environment=MEADOW_WD_INTERVAL=17
Environment=MEADOW_WD_JITTER=3
Environment=MEADOW_WD_WINDOW=180
Environment=MEADOW_WD_MIN_UNIQUE=2
Environment=MEADOW_WD_COOLDOWN=30
Environment=MEADOW_WD_RATE_WINDOW=3600
Environment=MEADOW_WD_MAX_RESTARTS=3

# IMPORTANT: run via bash to avoid exec permission/noexec issues
ExecStart=/bin/bash /home/meadow/meadow-kiosk/kiosk-freeze-watchdog.sh

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
