# systemd/meadow-vend-poller.service
[Unit]
Description=Meadow Vend Poller (WP queued vend -> local /vend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meadow
WorkingDirectory=/home/meadow/meadow-kiosk
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/meadow/meadow-kiosk/vend_poller.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
