# systemd/meadow-kiosk-browser.service
[Unit]
Description=Meadow Kiosk Browser (Chromium watchdog loop)
After=graphical.target network-online.target
Wants=network-online.target

StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=meadow
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/meadow/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

# Ensure kiosk is allowed to run
ExecStartPre=/bin/bash -lc 'rm -f /tmp/meadow_kiosk_stop || true'

ExecStart=/bin/bash /home/meadow/kiosk-browser.sh

ExecStop=/bin/bash -lc 'touch /tmp/meadow_kiosk_stop; pkill -f "chromium.*--kiosk" || true; pkill -f "chromium-browser.*--kiosk" || true'

Restart=on-failure
RestartSec=3

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
